# scripts/mhr-client-relay.service
[Unit]
Description=mhr-ggate client-side relay (links xray client to GAS)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=%h/.config/mhr-ggate
ExecStart=/usr/bin/env python3 %h/.config/mhr-ggate/client_relay.py --config %h/.config/mhr-ggate/relay.toml
Restart=on-failure
RestartSec=2

[Install]
WantedBy=default.target
